# .config/autostart/onedrive.desktop
[Desktop Entry]
Name=Onedrive
Exec=sh -c "rclone mount --daemon Onedrive: /mnt/Onedrive"
Terminal=false
Type=Application
